# templates/systemd.service.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

[Unit]
Description=Network daemon for aggregating statistics
Documentation=https://github.com/etsy/statsd/
Wants=network.target

[Service]
Type=simple
User={{statsd_user}}
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier={{statsd_title}}
ExecStart={{statsd_nodejs}} {{statsd_home}}/node_modules/statsd/stats.js {{statsd_home}}/config.js

[Install]
WantedBy=multi-user.target